# systemd/ota-worker@.service
[Unit]
Description=OTA scraper worker instance %i
After=network-online.target
StartLimitIntervalSec=3600
StartLimitBurst=5
[Service]
User=jf
Group=jf
WorkingDirectory=/home/jf/Projects/OTA-SCRAPER-LINUX
Environment=INSTANCE_ID=instance_%i
Environment=INSTANCE_DATA_DIR=/home/jf/Projects/OTA-SCRAPER-LINUX/data/instances/instance_%i
Environment=DB_BACKEND=sqlite
ExecStart=/usr/bin/systemd-inhibit --what=sleep:idle --who=OTA-SCRAPER --why=Long-running-OTA-scrape --mode=block /home/jf/Projects/OTA-SCRAPER-LINUX/.venv/bin/python tools/ota_worker_supervisor.py
Restart=on-failure
RestartSec=10
KillMode=control-group
TimeoutStopSec=45
MemoryAccounting=yes
CPUAccounting=yes
TasksAccounting=yes
TasksMax=256
